# ops/launchd/com.ida-trading.vps-tunnels.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ida-trading.vps-tunnels</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/ssh</string>
    <string>-i</string>
    <string>/Users/jenriquezafra/.ssh/trading_strats_hetzner_ed25519</string>
    <string>-N</string>
    <string>-L</string>
    <string>8700:127.0.0.1:8700</string>
    <string>-L</string>
    <string>5901:127.0.0.1:5901</string>
    <string>-o</string>
    <string>ExitOnForwardFailure=yes</string>
    <string>-o</string>
    <string>ServerAliveInterval=30</string>
    <string>-o</string>
    <string>ServerAliveCountMax=3</string>
    <string>trader@195.201.229.200</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/jenriquezafra/Library/Logs/ida-trading-vps-tunnels.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jenriquezafra/Library/Logs/ida-trading-vps-tunnels.err.log</string>
</dict>
</plist>
